# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/drf_tutorial/bin/uwsgi --ini /srv/subclass-drf/.config/uwsgi/uwsgi-app.ini